# cloudflared/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

xz zjkjzh fkjashdakd

# Declare variables
declare client_token

## Get the 'message' key from the user config options.
client_token=$(bashio::config 'token')

## Print the message the user supplied, defaults to "Hello World..."
bashio::log.info "Starting the cloudflare tunnel daemon..."

# Use the client_token variable in the cloudflared command
cloudflared tunnel --no-autoupdate run --token "${client_token}"

bashio::log.info "Cloudflare tunnel daemon started!"

# exec /usr/bin/my_program
