# cloudflared/rootfs/etc/services.d/cloudflared/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the addon service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Declare variables
declare client_token

## Get the 'message' key from the user config options.
client_token=$(bashio::config 'token')

## Print the message the user supplied, defaults to "Hello World..."
bashio::log.info "Starting the cloudflare tunnel daemon..."

# Use the client_token variable in the cloudflared command
exec cloudflared tunnel --no-autoupdate run --token "${client_token}"