# cloudflared/rootfs/etc/services.d/cloudflared/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the addon service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Declare variables
declare client_token

## Get the user config options.
client_token=$(bashio::config 'token')

bashio::log.info "Starting the cloudflare tunnel daemon..."

#exec /usr/bin/my_program
exec tunnel --no-autoupdate run --token "${client_token}" 1>&1 2>&1
# exec cloudflared tunnel --no-autoupdate run --token "${client_token}" 1>&1 2>&1
